# jetson/tractor-steering.service
[Unit]
Description=tractor steering service
ConditionPathExists=/tmp/tractor-ready.touch

[Service]
ExecStart=/home/farmer/tractor/jetson/steering.sh
WorkingDirectory=/home/farmer/tractor
StandardOutput=inherit
StandardError=inherit
Restart=always
User=farmer

[Install]
WantedBy=multi-user.target
