# jetson/tractor-camera.path
[Unit]
Description=Monitor /bin/tractor-ready.touch and start tractor-camera

[Path]
PathExists=/tmp/tractor-ready.touch
Unit=tractor-camera.service

[Install]
WantedBy=multi-user.target
