# jetson/tractor-programd.path
[Unit]
Description=Monitor /bin/tractor-ready.touch and start tractor-programd

[Path]
PathExists=/tmp/tractor-ready.touch
Unit=tractor-programd.service

[Install]
WantedBy=multi-user.target
