# jetson/tractor-steering.path
[Unit]
Description=Monitor /bin/tractor-ready.touch and start tractor-steering

[Path]
PathExists=/tmp/tractor-ready.touch
Unit=tractor-steering.service

[Install]
WantedBy=multi-user.target
